# deploy/backup.service
[Unit]
Description=Pick The Score Postgres backup → local disk
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/pckthscr
EnvironmentFile=/opt/pckthscr/.env
ExecStart=/opt/pckthscr/backup.sh
StandardOutput=journal
StandardError=journal
TimeoutStartSec=15min
